# worker/com.murmur.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.murmur.worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>/Users/samrobinson/Desktop/Murmur/murmur/worker/transcribe_worker.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MURMUR_API_URL</key>
        <string>https://192.168.1.214</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/murmur-worker.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/murmur-worker.log</string>
</dict>
</plist>
